# automation/launchd/agents/com.cosmdandy.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cosmdandy.backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/cosmdandy/.dotfiles/automation/backup/backup.sh</string>
    </array>
    <!-- Ежедневно в 13:00. Днём, а не ночью: ноутбук ночью спит, а launchd,
         в отличие от systemd с Persistent=true, пропущенный календарный
         запуск не догоняет. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>13</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <!-- Бэкап 5 ГБ не должен конкурировать за диск с интерактивной работой -->
    <key>LowPriorityIO</key>
    <true/>
    <key>Nice</key>
    <integer>5</integer>
    <key>StandardOutPath</key>
    <string>/tmp/backup.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/backup.log</string>
</dict>
</plist>
